# README

Included Files and locations :

/usr/local/bin/modemcomms.sh
/usr/local/bin/modemRegState.pl
/usr/local/bin/modemSigLevel.pl
/etc/network/interfaces.d/wwan0
/etc/sim.conf
/etc/qmi-network.conf


perl scripts require the following installed

apt-get install perl
cpan install Device::Modem
cpan install Device::Gsm 

//launch on boot
sudo cp /home/fledge/modemcomms/bin/modemcomms.sh /etc/init.d/
cd /etc/init.Device
sudo nano modemcomms.sh
///////add at top /////////
# /etc/init.d/modemcomms.sh
### BEGIN INIT INFO
# Provides:          modemcomms.sh
# Required-Start:    $remote_fs $syslog
# Required-Stop:     $remote_fs $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start daemon at boot time
# Description:       Enable service provided by daemon.
### END INIT INFO
////////////////////

sudo chmod +x modemcomms.sh
sudo update-rc.d modemcomms.sh defaults